# etc/init/deployer.conf
#!upstart
# http://howtonode.org/deploying-node-upstart-monit
description "automatically run node.js projects"
author      "Joachim kainz"

start on startup
stop on shutdown

script
    echo $$ > /var/run/deployer.pid
    exec sudo -u deployer /var/deployer/bin/deployer >> /var/log/deployer.log 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/deployer.log
end script

pre-stop script
    rm /var/run/deployer.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/deployer.log
end script
